# bins/redis/cluster/haproxy/docker/etc/supervisor.d/containerservices.ini
[supervisord]
nodaemon=false

[program:REPLACE_NODE_NAME]
command=haproxy -f /etc/haproxy.conf
process_name=%(program_name)s
directory=/opt/haproxy/node
user=root
autostart=true
autorestart=true
stopsignal=QUIT
redirect_stderr=true
stdout_logfile=/opt/data/dst/supervisornode.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stdout_capture_maxbytes=50MB

