# bins/redis/cluster/redis/docker/etc/supervisor.d/containerservices.ini
[supervisord]
nodaemon=false

[program:REPLACE_NODE_NAME]
command=redis-server /opt/redis/node/current.conf
process_name=%(program_name)s
directory=/opt/redis/node
user=root
autostart=true
autorestart=true
stopsignal=QUIT
redirect_stderr=true
stdout_logfile=/opt/data/dst/supervisor-redis.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stdout_capture_maxbytes=50MB

